# systemd/kgate-daemon.service
[Unit]
Description=Knowledge Gate Daemon
Documentation=https://github.com/loljeah/exambuilder
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/kgate-daemon
Restart=on-failure
RestartSec=5

# Socket path
Environment=XDG_RUNTIME_DIR=%t

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/kgate %h/.kgate %t
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=graphical-session.target
